# tools/goal-detector/com.filipeandrade.goal-detector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent for the WM goal detector.

  Install:
    security add-generic-password -a "$USER" -s wm-markers-token -w "<TOKEN>" -U
    cp tools/goal-detector/com.filipeandrade.goal-detector.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.filipeandrade.goal-detector.plist

  Logs:  tail -f ~/Library/Logs/wm-goal-detector.log
  State: ~/Developer/world-cup-viewer/tools/goal-detector/.state.json
-->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.filipeandrade.goal-detector</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>export WM_MARKERS_TOKEN="$(security find-generic-password -a $USER -s wm-markers-token -w)"; export WORKER_BASE="https://wm.filipeandrade.com"; exec /opt/homebrew/bin/node /Users/filipeandrade/Developer/world-cup-viewer/tools/goal-detector/detect.mjs</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>StartInterval</key><integer>1800</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/filipeandrade/Library/Logs/wm-goal-detector.log</string>
  <key>StandardErrorPath</key><string>/Users/filipeandrade/Library/Logs/wm-goal-detector.log</string>
  <key>WorkingDirectory</key><string>/Users/filipeandrade/Developer/world-cup-viewer/tools/goal-detector</string>
</dict>
</plist>
